# services/falcon-trader.service
[Unit]
Description=Falcon Trading Bot - Multi-Strategy Orchestrator
Documentation=https://github.com/TradingAsBuddies/falcon-trader
After=network-online.target falcon-screener@morning.service
Wants=network-online.target

[Service]
Type=simple
User=falcon
Group=falcon
WorkingDirectory=/var/lib/falcon

# Load environment
EnvironmentFile=/etc/falcon/secrets.env
Environment="PYTHONUNBUFFERED=1"
Environment="FALCON_ENV=production"

# Run the trader (pip-installed package)
ExecStart=/opt/falcon/venv/bin/falcon-trader

# Restart policy
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=falcon-trader

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/falcon /var/cache/falcon

# Resource limits (trading needs more resources)
MemoryMax=2G
CPUQuota=75%

[Install]
WantedBy=multi-user.target
