# services/falcon-screener@evening.timer
[Unit]
Description=Falcon Evening Screener (7:00 PM ET)

[Timer]
OnCalendar=*-*-* 19:00:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
